# trunk/packages/gstreamer/ti-video/files/etc/init.d/tivideo
#!/bin/sh /etc/rc.common
# Copyright (C) 2006 OpenWrt.org
START=45

# NOTE: MPEG4 not currently supported

GSTREAMER_BIN="/usr/bin/gst-launch-0.10"

# h264 encode/decode (accounts for using resizer if <740x480)
POOLS_h264=1x5250000,3x829440,3x691200

# MPEG4 (currently broken)
# mpeg4 encode/decode
#POOLS_mpeg4=1x2073600,2x1036800,3x691200,2x253440

# $1 - codec
load_modules() {
	local codec=$1
	local pools

	echo "Loading cmemk for $codec"
	eval pools=\$POOLS_${codec}
	echo "POOLS_$1=$pools"
	[ "$pools" ] || {
		echo "Unsupported codec: $codec"
		exit 1
	}
	insmod cmemk phys_start=0x87600000 phys_end=0x88000000 pools=$pools

	grep dsplinkk /proc/modules >/dev/null || {
		insmod dsplinkk
	}
}

load_video() {
	GSTREAMER_ARGS=""
	A_ARGS=""
	append A_ARGS "a=rtpmap:96"
	
	local type
	local host
	local port
	local fps
	local bps
	local codec
	local jitter
	local width
	local height
	local ratecontrol
	local quality
	local resolution
	
	config_get type $1 type "server"
	config_get host $1 host
	config_get port $1 port
	config_get fps $1 fps "30"
	config_get bps $1 bps "200000"
	config_get codec $1 codec "h264"
	config_get jitter $1 jitter "0"
	config_get resolution $1 resolution "720x480"
	config_get ratecontrol $1 ratecontrol "2"
	config_get quality $1 quality "3"
	width="${resolution%*x*}"
	height="${resolution##*x}"
	
	echo "setting type ${type}"
	echo "setting host ${host}"
	echo "setting port ${port}"
	echo "setting fps ${fps}"
	echo "setting bps ${bps}"
	echo "setting codec ${codec}"
	echo "setting resolution ${width}x${height}"
	
	grep cmemk /proc/modules >/dev/null || {
		load_modules $codec
	}

	[ "$type" == "server" ] && {
		[ "$fps" == "30" ] && {
			append GSTREAMER_ARGS "v4l2src always-copy=0 queue-size=1 !"	
		} || {
			append GSTREAMER_ARGS "v4l2src always-copy=1 queue-size=1 !"
			append GSTREAMER_ARGS "videorate ! video/x-raw-yuv,framerate=$fps/1 !"
		} 
		append GSTREAMER_ARGS "TIVidResize ! video/x-raw-yuv,width=$width,height=$height !"
		append GSTREAMER_ARGS "TIVidenc1 engineName=encode encodingPreset=$quality rateControlPreset=$ratecontrol bitRate=$bps frameRate=$fps"
		[ "$fps" == "30" ] && {
			append GSTREAMER_ARGS "contiguousInputFrame=1"
		}
		[ "$codec" == "h264" ] && {
			append A_ARGS "H264/90000"
			append GSTREAMER_ARGS "codecName=h264enc ! rtph264pay pt=96"
		} || {
			append A_ARGS "MP4V-ES/90000"
			append GSTREAMER_ARGS "codecName=mpeg4enc ! rtpmp4vpay pt=96"
		}
		append GSTREAMER_ARGS "! udpsink host=$host port=$port"
	} || {
		# TODO:
		#  - only use multicast-group/multicast-iface if host is mcast
		#  - add interface as a parameter

		# make sure mcast route exists or udpsrc multicast-group
		# causes gst-launch to exit with an error
		route add -net 224.0.0.0 netmask 240.0.0.0 br-lan
		append GSTREAMER_ARGS "udpsrc multicast-group=$host multicast-iface=br-lan port=$port caps=\"application/x-rtp, media=(string)video, clock-rate=(int)90000, profile-level-id=(string)42801e, payload=(int)96,"
		[ "$codec" == "h264" ] && {
			append GSTREAMER_ARGS "encoding-name=(string)H264\" !" 
		} || {
			append GSTREAMER_ARGS "encoding-name=(string)MP4V-ES\" !"
		}
		[ "$jitter" == "0" ] || {
			append GSTREAMER_ARGS "gstrtpjitterbuffer latency=$jitter !"
		}
		[ "$codec" == "h264" ] && {
			append GSTREAMER_ARGS "rtph264depay !"
		} || {
			append GSTREAMER_ARGS "rtpmp4vdepay !"
		} 
		append GSTREAMER_ARGS "TIViddec2 engineName=decode"
		[ "$codec" == "h264" ] && {
			append GSTREAMER_ARGS "codecName=h264dec !"
		} || {
			append GSTREAMER_ARGS "codecName=mpeg4dec !"
		}
		append GSTREAMER_ARGS "TIDmaiVideoSink videoStd=D1_NTSC videoOutput=composite resizer=1 sync=0"
	}
	
	echo "gstreamer args = $GSTREAMER_ARGS"
	start-stop-daemon -S -x $GSTREAMER_BIN \
		-p /var/run/video.pid \
		-m -b -- $GSTREAMER_ARGS
	
	rm -f /www/video.sdp
	echo "v=0" > /www/video.sdp
	echo "o=- 1188340656180883 1 IN IP4 $host/255" >> /www/video.sdp
	echo "s=Session streamed by GStreamer" >> /www/video.sdp
	echo "i=server.sh" >> /www/video.sdp
	echo "t=0 0" >> /www/video.sdp
	echo "a=tool:GStreamer" >> /www/video.sdp
	echo "a=type:broadcast" >> /www/video.sdp
	echo "m=video $port RTP/AVP 96" >> /www/video.sdp
	echo "c=IN IP4 $host/255" >> /www/video.sdp
	echo "$A_ARGS" >> /www/video.sdp	
}

start() {
	# Note: we dont support multiple video sections
	config_load tivideo
	config_foreach load_video video
}

stop() {
	[ -f /var/run/video.pid ] && {
		start-stop-daemon -K -q \
			-p /var/run/video.pid -s KILL
		
		rm -f /var/run/video.pid
		sleep 4
		# re-enable dymaic module load/unload when supporting multiple codecs
		#rmmod cmemk
	}
}
